# server01/etc/init/damphat.com.conf
#!upstart
description "webserver damphat.com"
author      "damphat"

start on startup
stop on shutdown

script
	cd /myprojects/damphat.com
	supervisor app.js
end script

